# crontab
*/30 * * * * root poetry --ansi --no-interaction run auto-tagger --add-tags deepbooru "type:image -deepbooru"
